# services/core4.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core4.sh

[Install]
WantedBy=default.target
