# services/core6.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core6.sh

[Install]
WantedBy=default.target
