# services/core5.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core5.sh

[Install]
WantedBy=default.target
