# services/core1.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core1.sh

[Install]
WantedBy=default.target
