# services/core2.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core2.sh

[Install]
WantedBy=default.target
