# services/core3.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/core3.sh

[Install]
WantedBy=default.target
