# services/mem2.service
[Unit]
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mem2.sh

[Install]
WantedBy=default.target
